# ops/schedule_backfill.cron
CRON_TZ=Asia/Riyadh
*/10 * * * * /usr/local/bin/python /app/ops/train_backfill_one.py >> /app/logs/cron/backfill.log 2>&1
